# .config/autostart/bing-wallpaper.desktop
[Desktop Entry]
Type=Application
Name=Bing Wallpaper
Exec=/home/jesse/Pictures/get_bing_wallpaper.py
Path=/home/jesse/Pictures
Terminal=false
